# deploy/fitness-slack-bot.timer
[Unit]
Description=Daily fitness Slack post (weekdays at 08:55 America/New_York)

[Timer]
OnCalendar=Mon..Fri 08:55 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
